# rclone-bisync.service
[Unit]
Description=Rclone Bisync Daemon
After=network.target

[Service]
ExecStart=/usr/bin/rclone-bisync --daemon
Restart=on-failure

[Install]
WantedBy=multi-user.target